# .config/autostart/dropbox.desktop
[Desktop Entry]
Name=dropbox
Exec=run_dropbox
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
