# .config/autostart/dropbox.desktop
[Desktop Entry]
Name=dropbox
Exec=run_dropbox.sh
Terminal=false
Type=Application
