# .config/autostart/flameshot.desktop
[Desktop Entry]
Name=flameshot
Exec=flameshot
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=false
